# launchagents/com.localllm.server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.localllm.server</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/jrg/local-llm/venv/bin/mlx-omni-server</string>
        <string>--host</string>
        <string>127.0.0.1</string>
        <string>--port</string>
        <string>10240</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
        <key>VIRTUAL_ENV</key>
        <string>/Users/jrg/local-llm/venv</string>
        <key>PATH</key>
        <string>/Users/jrg/local-llm/venv/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HF_HUB_OFFLINE</key>
        <string>1</string>
        <key>TRANSFORMERS_OFFLINE</key>
        <string>1</string>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <false/>
    <key>WorkingDirectory</key>
    <string>/Users/jrg/local-llm</string>
    <key>StandardOutPath</key>
    <string>/Users/jrg/local-llm/logs/llm.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/jrg/local-llm/logs/llm-error.log</string>
</dict>
</plist>
