# launchagents/com.belvedere.llm-server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.belvedere.llm-server</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/jrg/tts/.venv/bin/python</string>
        <string>-m</string>
        <string>mlx_lm.server</string>
        <string>--model</string>
        <string>Jiunsong/supergemma4-26b-uncensored-mlx-4bit-v2</string>
        <string>--port</string>
        <string>8800</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/tmp/belvedere-llm.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/belvedere-llm.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
